# Code - OSS/User/History/46782be1/Dh7L.service
[Unit]
Description=Set Hyprsunset temperature at boot

[Service]
Type=oneshot
ExecStart=/usr/bin/set-hyprsunset.sh

[Install]
WantedBy=multi-user.target
